# .config/systemd/user/podentr.service
[Unit]
Description=Automate podboat

[Service]
Type=simple
ExecStart=%h/.local/bin/podentr

[Install]
WantedBy=default.target
